# deploy/ecolyxis.service.d/secret-key.conf
# Drop-in for ecolyxis.service — fixes the weak in-repo SECRET_KEY.
#
# The main unit injects `Environment=SECRET_KEY=ecolyxis-prod-change-this-later`,
# which (because config.py loads .env with setdefault — real env wins) overrides
# the strong key in /opt/Ecolyxis/.env, so Flask sessions are signed with a
# guessable, version-controlled value.
#
# UnsetEnvironment is applied AFTER the main unit's Environment=, removing
# SECRET_KEY from the process env. config.py then sources the strong key from
# /opt/Ecolyxis/.env instead. Verified: app resolves the 37-char .env key.
#
# Install (root):
#   sudo cp deploy/ecolyxis.service.d/secret-key.conf /etc/systemd/system/ecolyxis.service.d/
#   sudo systemctl daemon-reload && sudo systemctl restart ecolyxis
# Verify (expect 0):
#   sudo tr '\0' '\n' < /proc/$(systemctl show -p MainPID --value ecolyxis)/environ | grep -c '^SECRET_KEY='
# Revert: delete this file + daemon-reload + restart.
#
# NOTE: switching keys invalidates existing session cookies once (all users
# logged out a single time) — do it during the cutover restart window.
[Service]
UnsetEnvironment=SECRET_KEY
